# supervisor-consumers.ini
[supervisord]
user=root

[program:nginx]
command=nginx -g "daemon off;"
process_name=%(program_name)s
autostart=true
autorestart=unexpected

[program:php-fpm]
command=php-fpm
process_name=%(program_name)s
autostart=true
autorestart=unexpected